# services.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
logfile_maxbytes=50MB
logfile_backups=5
loglevel=info

[program:redis]
command=/usr/bin/redis-server /etc/redis/redis.conf --port 6380 --requirepass trading_system_2025 --dir /data
autostart=true
autorestart=true
stderr_logfile=/var/log/redis/redis-server.err.log
stdout_logfile=/var/log/redis/redis-server.out.log
priority=10
startsecs=10
startretries=5
user=root
stopasgroup=true
killasgroup=true
stopsignal=TERM
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB

[program:prometheus]
command=/usr/bin/prometheus --config.file=/etc/prometheus/prometheus.yml --storage.tsdb.path=/prometheus --web.listen-address=:9090 --web.enable-lifecycle --web.enable-admin-api
autostart=true
autorestart=true
stderr_logfile=/var/log/prometheus.err.log
stdout_logfile=/var/log/prometheus.out.log
priority=20
user=root
environment=PATH="/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
directory=/prometheus
startretries=5
startsecs=3

[program:redis_exporter]
command=/usr/local/bin/redis_exporter --redis.addr=redis://localhost:6380 --redis.password=trading_system_2025
autostart=true
autorestart=true
stderr_logfile=/var/log/redis_exporter.err.log
stdout_logfile=/var/log/redis_exporter.out.log
priority=30
startsecs=3

[program:jupyter]
command=jupyter lab --ip=0.0.0.0 --port=8888 --no-browser --allow-root --NotebookApp.token=""
autostart=true
autorestart=true
stderr_logfile=/var/log/jupyter.err.log
stdout_logfile=/var/log/jupyter.out.log
priority=40

[program:trading_system]
command=python3 -u /app/project/scripts/start_system.py --health-check-interval=300 --debug
directory=/app/project
autostart=true
autorestart=true
stderr_logfile=/var/log/trading_system.err.log
stdout_logfile=/var/log/trading_system.out.log
priority=50
environment=PYTHONPATH="/app/project",REDIS_HOST="localhost",REDIS_PORT="6380",REDIS_PASSWORD="trading_system_2025",USE_REDIS_FALLBACK="true",SKIP_DEPENDENCIES_CHECK="false",DEBUG_MODE="true",PYTHONUNBUFFERED="1",NVIDIA_VISIBLE_DEVICES="all",CUDA_VISIBLE_DEVICES="all",TF_FORCE_GPU_ALLOW_GROWTH="true",TF_FORCE_UNIFIED_MEMORY="1"
stopwaitsecs=60
stopsignal=TERM
stopasgroup=true
killasgroup=true
startsecs=10
startretries=5

[program:frontend]
command=python -u /app/project/frontend/app.py
directory=/app/project
autostart=true
autorestart=true
stderr_logfile=/var/log/frontend.err.log
stdout_logfile=/var/log/frontend.out.log
priority=60
environment=PYTHONPATH="/app/project",FLASK_ENV="development",FLASK_DEBUG="1",FLASK_RUN_PORT="5000",FLASK_RUN_HOST="0.0.0.0",FLASK_APP="/app/project/frontend/app.py",FLASK_RUN_EXTRA_FILES="/app/project/frontend/templates/*",FRONTEND_WEBSOCKET_ENABLED="true",FRONTEND_REALTIME_UPDATES="true",PYTHONUNBUFFERED="1",REDIS_HOST="localhost",REDIS_PORT="6380",REDIS_PASSWORD="trading_system_2025",NVIDIA_VISIBLE_DEVICES="all",CUDA_VISIBLE_DEVICES="all"
stopwaitsecs=30
stopsignal=TERM
stopasgroup=true
killasgroup=true
startretries=10
startsecs=5

[program:redis_event_listener]
command=python -m frontend.event_listener
directory=/app/project
autostart=true
autorestart=true
stderr_logfile=/var/log/redis_event_listener.err.log
stdout_logfile=/var/log/redis_event_listener.out.log
priority=65
environment=PYTHONPATH="/app/project",REDIS_PUBSUB_THREADS="4",REDIS_NOTIFY_KEYSPACE_EVENTS="Kxe"
stopwaitsecs=20
stopsignal=TERM
stopasgroup=true
killasgroup=true

[program:monitoring_system]
command=python -m monitoring_system.monitoring_system
directory=/app/project
autostart=true
autorestart=true
stderr_logfile=/var/log/monitoring_system.err.log
stdout_logfile=/var/log/monitoring_system.out.log
priority=55
environment=PYTHONPATH="/app/project",REDIS_HOST="localhost",REDIS_PORT="6380",REDIS_PASSWORD="trading_system_2025",METRICS_PORT="8000",METRICS_COLLECTION_INTERVAL="15",PROMETHEUS_ENABLED="true",PYTHONUNBUFFERED="1",NVIDIA_VISIBLE_DEVICES="all",CUDA_VISIBLE_DEVICES="all"
stopwaitsecs=30
stopsignal=TERM
stopasgroup=true
killasgroup=true

[program:portfolio_updater]
command=python -u /app/project/update_portfolio.py
directory=/app/project
autostart=true
autorestart=true
stderr_logfile=/var/log/portfolio_updater.err.log
stdout_logfile=/var/log/portfolio_updater.out.log
priority=70
environment=PYTHONPATH="/app/project",REDIS_HOST="localhost",REDIS_PORT="6380",REDIS_PASSWORD="trading_system_2025",ALPACA_API_KEY="%(ENV_ALPACA_API_KEY)s",ALPACA_API_SECRET="%(ENV_ALPACA_API_SECRET)s",ALPACA_BASE_URL="%(ENV_ALPACA_API_URL)s",PYTHONUNBUFFERED="1"
stopwaitsecs=20
stopsignal=TERM
stopasgroup=true
killasgroup=true
startsecs=5
startretries=3
